# infra/vault-autocommit.service
[Unit]
Description=Auto-commit vault changes to git

[Service]
Type=oneshot
# %h expands to $HOME. Change "Projects/vault" if your vault lives elsewhere.
WorkingDirectory=%h/Projects/vault
ExecStart=/bin/bash -c '\
  git add -A && \
  if ! git diff --cached --quiet; then \
    git commit -m "auto: $(date +%%FT%%T)"; \
    git push --quiet || echo "[autocommit] push failed; next sync will retry" >&2; \
  fi'
