# self-hosting/vault-autocommit.service
[Unit]
Description=Auto-commit vault changes to git

[Service]
Type=oneshot
# %h expands to $HOME. Change "vault" if your vault lives elsewhere.
WorkingDirectory=%h/vault
ExecStart=/bin/bash -c '\
  git add -A && \
  if ! git diff --cached --quiet; then \
    git commit -m "auto: $(date +%%FT%%T)"; \
    git push --quiet 2>/dev/null || true; \
  fi'
